# systemd/closed.target
# This file is part of loom,
# hacker^Wspace automation kit

[Unit]
Description=Indispensible Loom Systems
